# debian/fbkeyboard.init
#!/sbin/openrc-run

description="fbkeyboard"
name="fbkeyboard"

start() {
	font=$(fc-match --format=%{file} DejaVuSans.tff)
	start-stop-daemon -S -b -x /usr/bin/fbkeyboard -- $font
}

stop() {
	start-stop-daemon -K --oknodo -x /usr/bin/fbkeyboard
}

