# debian/fbkeyboard.postinst
#!/bin/sh

install -d /etc/runlevels/shell/
rc-update del fbkeyboard default || true
rc-update add fbkeyboard shell   || true
