# debian/fbkeyboard.postinst
#!/bin/sh

install -d /etc/runlevels/shell/
rc-update del fbkeyboard default
rc-update add fbkeyboard shell
